# net/netbird/files/netbird.init
#!/bin/sh /etc/rc.common

START=99
STOP=10

USE_PROCD=1

start_service() {
	procd_open_instance
	procd_set_param command /usr/bin/netbird
	procd_append_param command service run
	procd_set_param pidfile /var/run/netbird.pid
	procd_close_instance
}
