# admin/muninlite/files/etc/xinetd.d/munin
service munin
{
	socket_type	= stream
	protocol	= tcp
	wait		= no
	user		= root
	group		= root
	server		= /usr/sbin/munin-node
	disable		= no
}
