# net/knxd/files/knxd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=98
STOP=20
NAME=knxd
PROG=/usr/bin/$NAME
USE_PROCD=1

start_service() {
	local options url
	config_load "$NAME"
	config_get options args options ''
	config_get url args url
	procd_open_instance
	procd_set_param command $PROG $options $url
	procd_set_param respawn
	procd_close_instance
}

stop_service() {
	killall knxd
}

reload_service()
{
	restart
}

service_triggers()
{
	procd_add_reload_trigger knxd
}
