# net/knxd/files/knxd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=98
STOP=20
NAME=knxd
PROG=/usr/bin/$NAME

start() {
        local options url
        config_load "$NAME"
        config_get options args options ''
        config_get url args url
        service_start $PROG $options $url
}

stop() {
        service_stop $PROG
}
