# net/netopeer2/files/netopeer2-server.init
#!/bin/sh /etc/rc.common

START=90
STOP=10

USE_PROCD=1
PROG=/bin/netopeer2-server

start_service() {

    # netopeer2-server requires sysrepo daemon
    /etc/init.d/sysrepo start

    sleep 1

    procd_open_instance
    procd_set_param command $PROG
    procd_append_param command -d -v 0
    procd_set_param respawn
    procd_close_instance
}

stop_service()
{
    service_stop ${PROG}
}
