# utils/rng-tools/files/rngd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2011 OpenWrt.org

START=98

RNGD_INTERVAL=30
RNGD_AMOUNT=4000
RNGD_DEVICE="/dev/urandom"

start() {
	service_start /sbin/rngd -r $RNGD_DEVICE -W $RNGD_AMOUNT -t $RNGD_INTERVAL
}

stop() {
	service_stop /sbin/rngd
}
