# net/ddns-scripts/files/etc/hotplug.d/iface/25-ddns
#!/bin/sh

if [ "$ACTION" = "ifup" ]; then
	. /usr/lib/ddns/dynamic_dns_functions.sh
	/etc/init.d/ddns enabled && start_daemon_for_all_ddns_sections "$INTERFACE"
fi

