# net/ddns-scripts/files/etc/hotplug.d/iface/25-ddns
#!/bin/sh

. /usr/lib/ddns/dynamic_dns_functions.sh

if [ "$ACTION" = "ifup" ]; then
	start_daemon_for_all_ddns_sections "$INTERFACE"
fi


