# net/ddns-scripts/files/etc/hotplug.d/iface/ddns
#!/bin/sh

# there are other ACTIONs we don't need
case "$ACTION" in
	ifup)					# OpenWrt is giving a network not phys. Interface
		/etc/init.d/ddns enabled && /usr/lib/ddns/dynamic_dns_updater.sh -n "$INTERFACE" -- start
		;;
	ifdown)
		/usr/lib/ddns/dynamic_dns_updater.sh -n "$INTERFACE" -- stop
		;;
	ifupdate)
		[ -z "$IFUPDATE_ADDRESSES" ] && exit 0
		/usr/lib/ddns/dynamic_dns_updater.sh -n "$INTERFACE" -- stop
		sleep 1	# give time to shutdown
		/etc/init.d/ddns enabled && /usr/lib/ddns/dynamic_dns_updater.sh -n "$INTERFACE" -- start
		;;
esac
