# net/ddns-scripts/files/etc/hotplug.d/iface/ddns
#!/bin/sh

start_updater() {
	/usr/lib/ddns/dynamic_dns_updater.sh -n "$INTERFACE" -- start
	/usr/lib/ddns/dynamic_dns_updater.sh -n "$DEVICE" -- start
}

stop_updater() {
	/usr/lib/ddns/dynamic_dns_updater.sh -n "$INTERFACE" -- stop
	/usr/lib/ddns/dynamic_dns_updater.sh -n "$DEVICE" -- stop
}

# there are other ACTIONs like ifupdate we don't need
case "$ACTION" in
	ifup)
		/etc/init.d/ddns enabled && start_updater
		;;
	ifdown)
		stop_updater
		;;
esac
