# net/mwan3/files/etc/hotplug.d/iface/15-mwan3
#!/bin/sh

. /lib/functions.sh
. /lib/functions/network.sh
. /lib/mwan3/mwan3.sh
. /lib/mwan3/common.sh

initscript=/etc/init.d/mwan3
. /lib/functions/procd.sh


SCRIPTNAME="mwan3-hotplug"
[ "$ACTION" = "ifup" ] || [ "$ACTION" = "ifdown" ] || [ "$ACTION" = "connected" ] || [ "$ACTION" = "disconnected" ] || exit 1
[ -n "$INTERFACE" ] || exit 2
[ "$FIRSTCONNECT" = "1" ] || [ "$MWAN3_SHUTDOWN" = "1" ] && exit 0

if { [ "$ACTION" = "ifup" ] || [ "$ACTION" = "connected" ] ; } && [ -z "$DEVICE" ]; then
	LOG notice "$ACTION called on $INTERFACE with no device set"
	exit 3
fi

[ "$MWAN3_STARTUP" = "init" ] || procd_lock

mwan3_init

/etc/init.d/mwan3 running || {
	[ "$MWAN3_STARTUP" = "init" ] || procd_lock
	LOG notice "mwan3 hotplug $ACTION on $INTERFACE not called because globally disabled"
	mwan3_flush_conntrack "$INTERFACE" "$ACTION"
	exit 0
}

$IPT4 -S mwan3_hook &>/dev/null || {
	LOG warn "hotplug called on $INTERFACE before mwan3 has been set up"
	exit 0
}

if [ "$MWAN3_STARTUP" != "init" ] && [ "$ACTION" = "ifup" ]; then
	mwan3_set_user_iface_rules $INTERFACE $DEVICE
fi

config_get_bool enabled $INTERFACE 'enabled' '0'
[ "${enabled}" -eq 1 ] || {
	LOG notice "mwan3 hotplug on $INTERFACE not called because interface disabled"
	exit 0
}

config_get initial_state $INTERFACE initial_state "online"
if [ "$initial_state" = "offline" ]; then
	status=$(cat $MWAN3TRACK_STATUS_DIR/$INTERFACE/STATUS 2>/dev/null || echo unknown)
	[ "$status" = "online" ] || status=offline
else
	status=online
fi

LOG notice "Execute $ACTION event on interface $INTERFACE (${DEVICE:-unknown})"

case "$ACTION" in
	connected)
		mwan3_set_iface_hotplug_state $INTERFACE "online"
		mwan3_set_policies_iptables
		;;
	ifup)
		mwan3_create_iface_iptables $INTERFACE $DEVICE
		mwan3_create_iface_rules $INTERFACE $DEVICE
		mwan3_set_iface_hotplug_state $INTERFACE "$status"
		if [ "$MWAN3_STARTUP" != "init" ]; then
			mwan3_create_iface_route $INTERFACE $DEVICE
			mwan3_set_general_rules
			[ "$status" = "online" ] && mwan3_set_policies_iptables
		fi
		[ "$ACTION" = ifup ] && procd_running mwan3 "track_$INTERFACE" && procd_send_signal mwan3 "track_$INTERFACE" USR2
		;;
	disconnected)
		mwan3_set_iface_hotplug_state $INTERFACE "offline"
		mwan3_set_policies_iptables
		;;
	ifdown)
		mwan3_set_iface_hotplug_state $INTERFACE "offline"
		mwan3_delete_iface_ipset_entries $INTERFACE
		mwan3_delete_iface_rules $INTERFACE
		mwan3_delete_iface_route $INTERFACE
		mwan3_delete_iface_iptables $INTERFACE
		procd_running mwan3 "track_$INTERFACE" && procd_send_signal mwan3 "track_$INTERFACE" USR1
		mwan3_set_policies_iptables
	;;
esac
exit 0
