# net/mwan3/files/etc/hotplug.d/iface/15-mwan3
#!/bin/sh

. /lib/functions.sh
. /lib/functions/network.sh
. /lib/mwan3/mwan3.sh

[ "$ACTION" == "ifup" -o "$ACTION" == "ifdown" ] || exit 1
[ -n "$INTERFACE" ] || exit 2

if [ "$ACTION" == "ifup" ]; then
        [ -n "$DEVICE" ] || exit 3
fi

[ -x /usr/bin/ip ] || exit 4
[ -x /usr/sbin/ipset ] || exit 5
[ -x /usr/sbin/iptables ] || exit 6
[ -x /usr/sbin/ip6tables ] || exit 7
[ -x /usr/bin/logger ] || exit 8

mwan3_lock
mwan3_set_connected_iptables
mwan3_unlock

config_load mwan3
config_get enabled $INTERFACE enabled 0
[ "$enabled" == "1" ] || exit 0

config_get family $INTERFACE family ipv4

if [ "$family" == "ipv4" ]; then
	network_get_gateway gateway $INTERFACE
elif [ "$family" == "ipv6" ]; then
	network_get_gateway6 gateway $INTERFACE
fi

if [ "$ACTION" == "ifup" ]; then
	[ -n "$gateway" ] || exit 9
fi

mwan3_lock
$LOG notice "$ACTION interface $INTERFACE (${DEVICE:-unknown})"

case "$ACTION" in
	ifup)
		mwan3_set_general_rules
		mwan3_set_general_iptables
		mwan3_create_iface_rules $INTERFACE $DEVICE
		mwan3_create_iface_iptables $INTERFACE $DEVICE
		mwan3_create_iface_route $INTERFACE $DEVICE
		mwan3_set_iface_hotplug_state $INTERFACE "online"
		mwan3_track $INTERFACE $DEVICE
		mwan3_set_policies_iptables
		mwan3_set_user_rules
	;;
	ifdown)
		mwan3_set_iface_hotplug_state $INTERFACE "offline"
		mwan3_delete_iface_ipset_entries $INTERFACE
		mwan3_set_policies_iptables
		mwan3_set_user_rules
	;;
esac

mwan3_unlock

exit 0
