# net/mwan3/files/etc/hotplug.d/iface/16-mwan3-user
#!/bin/sh

[ -f "/etc/mwan3.user" ] && {
	. /lib/functions.sh
	. /lib/mwan3/mwan3.sh

	mwan3_lock "$ACTION" "user"

	config_load mwan3
	config_get_bool enabled globals 'enabled' '0'
	[ "${enabled}" -gt 0 ] || {
		mwan3_unlock "$ACTION" "user"
		exit 0
	}

	config_get_bool enabled "$INTERFACE" enabled 0
	[ "${enabled}" -eq 1 ] || {
		mwan3_unlock "$ACTION" "user"
		exit 0
	}

	mwan3_unlock "$ACTION" "user"

	env -i ACTION="$ACTION" INTERFACE="$INTERFACE" DEVICE="$DEVICE" \
		/bin/sh /etc/mwan3.user
}

exit 0
