# net/mwan3/files/etc/hotplug.d/iface/16-mwan3-user
#!/bin/sh

[ -f "/etc/mwan3.user" ] && {
	. /lib/functions.sh

	config_load mwan3
	config_get enabled $INTERFACE enabled 0
	[ "${enabled}" = "1" ] || exit 0
	/bin/sh /etc/mwan3.user $ACTION $INTERFACE $DEVICE
}

exit 0
