# net/c-icap-server/files/c-icap.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2013 OpenWrt.org

START=20
USE_PROCD=1

start_service()
{
        procd_open_instance
        procd_set_param command "/usr/bin/c-icap" -N 
        procd_close_instance
}

stop_service()
{
        kill `cat /var/run/c-icap.pid`
}

