# net/noip/files/etc/init.d/noip
#!/bin/sh /etc/rc.common
START=99

BIN=noip2

start() {
    $BIN 
}

stop() {
    $BIN -K `pgrep $BIN`
}

reload() {
    stop
    start         
}
