# utils/dvmhost/files/dvmhost.init
#!/bin/sh /etc/rc.common

START=80
USE_PROCD=1
PROG=/usr/bin/dvmhost
CONFIGFILE=/etc/dvmhost/config.yml

start_service() {
	procd_open_instance
	procd_set_param command $PROG -c
	procd_set_param file $CONFIGFILE
	procd_close_instance
}
