# utils/telegraf/files/etc/uci-defaults/90-telegraf
#!/bin/sh
# Generates default telegraf.conf if it doesn't already exist

[ -f /etc/telegraf.conf ] && exit 0

/usr/bin/telegraf config create > /etc/telegraf.conf

exit 0
