# net/https-dns-proxy/files/etc/hotplug.d/online/30-https-dns-proxy
#!/bin/sh
/etc/init.d/https-dns-proxy start 'on_hotplug'
