# net/acme-common/files/acme.init
#!/bin/sh /etc/rc.common

START=80
USE_PROCD=1

service_triggers() {
	procd_add_config_trigger config.change acme \
		/usr/bin/acme get
}
