# net/fcgiwrap/files/fcgiwrap.init
#!/bin/sh /etc/rc.common

START=60

USE_PROCD=1

start_service() {
	procd_open_instance
	procd_set_param command /usr/bin/spawn-fcgi -n -s /var/run/cgi.sock -- /usr/bin/fcgiwrap
	procd_close_instance
}
