# net/mptcpd/files/mptcpd.init
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=10

start_service() {
	procd_open_instance
	procd_set_param command /usr/bin/mptcpd
	procd_set_param respawn
	procd_close_instance
}
