# net/keepalived/files/etc/hotplug.d/keepalived/509-ucitrack
#!/bin/sh

# shellcheck source=/dev/null
. /lib/functions/keepalived/hotplug.sh

set_service_name ucitrack

set_reload_if_sync

add_sync_file /etc/config/ucitrack

keepalived_hotplug
