# admin/sudo/files/sudo.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2014 OpenWrt.org

START=99

start() {
	[ -d /var/db/sudo ] || {
		mkdir -m 0755 -p /var/db/sudo
		chmod 0700 /var/db/sudo
	}
}
