# net/wgsd/files/wgsd-coredns.init
#!/bin/sh /etc/rc.common

START=99

USE_PROCD=1
PROG=/usr/bin/wgsd-coredns

start_service() {
	procd_open_instance
	procd_set_param command "$PROG" -conf /etc/Corefile
	procd_close_instance
}
