# net/zerotier/files/etc/hotplug.d/net/30-zerotier
#!/bin/sh

[ "$ACTION" = "add" ] || exit 0
echo "$DEVICENAME" | grep -Eq "^zt" || exit 0

[ "$(uci -q get zerotier.global.enabled)" = "1" ] || exit 0
zerotier-fw4 -i "$DEVICENAME"
