# net/nginx-utils/src/crontabs
3 3 12 12 * /etc/init.d/nginx 'add_ssl' '_lan'
