# net/bird/files/bird4.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2010-2011 OpenWrt.org

START=50

SERVICE_DAEMONIZE=1
SERVICE_WRITE_PID=1

start() {
	service_start /usr/sbin/bird4 -d
#	( SERVICE_MATCH_NAME=1 service_start /usr/sbin/bird4loop )
}

stop() {
#	( SERVICE_MATCH_NAME=1 service_stop /usr/sbin/bird4loop )
	service_stop /usr/sbin/bird4
}

reload() {
	service_reload /usr/sbin/bird4
}
