# net/arpwatch/files/arpwatch.init
#!/bin/sh /etc/rc.common
START=90

start () {
	arpwatch -f /etc/arpwatch/arp.dat -i br-lan
}

stop() {
	killall arpwatch
}
