# net/ddns-scripts/files/etc/init.d/ddns
#!/bin/sh /etc/rc.common
START=95

start() {
	. /usr/lib/ddns/dynamic_dns_functions.sh
	start_daemon_for_all_ddns_sections
}

stop() {
	killall -9 dynamic_dns_updater.sh
}

