# net/sqm-scripts/files/etc/init.d/sqm
#!/bin/sh /etc/rc.common

START=50

reload()
{
/usr/lib/sqm/run.sh start
}

restart()
{
reload
}

start()
{
reload
}

stop()
{
/usr/lib/sqm/run.sh stop
}