# oonf-dlep-proxy/files/dlep_proxy.init
#!/bin/sh /etc/rc.common

START=82
DAEMON='dlep_proxy'

. /lib/functions/oonf_init.sh
