# utils/btrfs-progs/files/btrfs-scan.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2014 OpenWrt.org

START=19

start() {
	grep -q btrfs /proc/filesystems && /usr/bin/btrfs device scan
}

