# net/sqm-scripts/files/etc/hotplug.d/iface/11-sqm
#!/bin/sh

[ "$ACTION" = ifup ] && /etc/init.d/sqm enabled && /usr/lib/sqm/run.sh start ${DEVICE}
