# files/etc/init.d/batmand
#!/bin/sh /etc/rc.common
START=90
. /lib/config/uci.sh
uci_load batmand
start () {
	interface=$(uci get batmand.general.interface)
	if [ "$interface" = "" ]; then
		echo $1 Error, you must specify at least a network interface
		exit
	fi
	announce=$(uci get batmand.general.announce)
	gateway_class=$(uci get batmand.general.gateway_class)
	originator_interval=$(uci get batmand.general.originator_interval)
	preferred_gateway=$(uci get batmand.general.preferred_gateway)
	routing_class=$(uci get batmand.general.routing_class)
	visualisation_srv=$(uci get batmand.general.visualisation_srv)
	policy_routing_script=$(uci get batmand.general.policy_routing_script)
	batman_args=""

	if [ $announce ]; then
		batman_args=${batman_args}'-a '$announce' '
	fi

	if [ $gateway_class ]; then
		batman_args=${batman_args}'-g '$gateway_class' '
	fi

	if [ $originator_interval ]; then
		batman_args=${batman_args}'-o '$originator_interval' '
	fi

	if [ $preferred_gateway ]; then
		batman_args=${batman_args}'-p '$preferred_gateway' '
	fi

	if [ $routing_class ]; then
		batman_args=${batman_args}'-r '$routing_class' '
	fi
		
	if [ $visualisation_srv ]; then
		batman_args=${batman_args}'-s '$visualisation_srv' '
	fi
	
	if [ $policy_routing_script ]; then
		batman_args=${batman_args}'--policy-routing-script '$policy_routing_script' '
	fi

	batman_args=${batman_args}$interface
	batmand $batman_args >/dev/null 2>&1
}

restart () {
	$0 stop
	sleep 3
	$0 start
}
        
stop () {
	killall batmand
}
