# files/etc/init.d/batman-adv
#!/bin/sh /etc/rc.common
# Copyright (C) 2011 OpenWrt.org

START=90

. /lib/batman-adv/config.sh

# can also be used with "batman-adv start bat0"
start() {
	config_load batman-adv

	if [ -n "$1" ]; then
		start_mesh $1
	else
		config_foreach start_mesh mesh
	fi
}

# can also be used with "batman-adv stop bat0"
stop () {
	config_load batman-adv

	if [ -n "$1" ]; then
		stop_mesh $1
	else
		config_foreach stop_mesh mesh
	fi
}
