# net/wifidog-ng/files/wifidog-ng.init
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=95

BIN=/usr/bin/wifidog-ng

start_service() {
    modprobe wifidog-ng

    procd_open_instance
    procd_set_param command $BIN
    procd_set_param respawn
    procd_close_instance
}

stop_service() {
    rmmod wifidog-ng
}
