# utils/prometheus-node-exporter-lua/files/etc/init.d/prometheus-node-exporter-lua
#!/bin/sh /etc/rc.common
# Copyright (C) 2013-2017 OpenWrt.org

START=60
USE_PROCD=1

_log() {
	logger -p daemon.info -t prometheus-node-exporter-lua "$@"
}

start_service() {
	. /lib/functions/network.sh

	local interface family ipv6 port bind bind6 bind4

	config_load prometheus-node-exporter-lua.main
	config_get interface "main" listen_interface "loopback"
	# Compatibility with old conf
	config_get_bool ipv6 "main" listen_ipv6
	config_get family "main" listen_family "any"
	config_get port "main" listen_port 9100

	# Migration from ipv6 to family
	if [ "$ipv6" ]; then
		_log "please migrate from ipv6 to family"
		case "$ipv6" in
			0) family="ipv4";;
			1) family="ipv6";;
		esac
		unset ipv6
	fi

	if [ "$interface" != "*" ]; then
		if [ "$family" = "any" ] || [ "$family" = "ipv4" ]; then
			network_get_ipaddrs bind4 "$interface"
		fi
		if [ "$family" = "any" ] || [ "$family" = "ipv6" ]; then
			network_get_ipaddrs6 bind6 "$interface"
		fi
		network_is_up "$interface" && ([ -n "$bind6" ] || [ -n "$bind4" ]) || {
			_log "defering start until listen interface $interface becomes ready"
			return 0
		}
		set -- $bind4 $bind6
	else
		set -- ::
	fi

	for bind; do
		procd_open_instance "$bind"

		procd_set_param command /usr/bin/prometheus-node-exporter-lua
		procd_append_param command --bind ${bind}
		procd_append_param command --port ${port}
		case "$family" in
			"ipv4") procd_append_param command -4;;
			"ipv6") procd_append_param command -6;;
		esac

		procd_set_param stdout 1
		procd_set_param stderr 1
		procd_set_param respawn

		procd_close_instance
	done
}

service_triggers()
{
	local interface

	procd_add_reload_trigger "prometheus-node-exporter-lua"

	config_load prometheus-node-exporter-lua.main
	config_get interface "main" listen_interface "loopback"

	[ "$interface" = "*" ] || procd_add_reload_interface_trigger "$interface"
}
