# net/simple-captive-portal/files/etc/hotplug.d/net/00-simple-captive-portal
INTF=$(uci -q get simple-captive-portal.main.interface)

if [ "$ACTION" = add -a "$DEVICENAME" == "$INTF" ]; then
    /etc/init.d/simple-captive-portal firewall
fi
