# utils/gpsd/files/etc/hotplug.d/gpsd/00-debug
#!/bin/sh

gps_log() {
	local level="$1"; shift

	[ "$level" = "debug" ] && return
	logger -p "daemon.${level}" -t "gpsd(debug)[$$]" "hotplug: $*"
}

main() {
	gps_log info "$(env)"
}

main "$@"
