# net/rp-pppoe/files/pppoe-relay.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=50

DEFAULT=/etc/default/pppoe-relay

start() {
	[ -f $DEFAULT ] && . $DEFAULT
	service_start /usr/sbin/pppoe-relay $OPTIONS
}

stop() {
	service_stop /usr/sbin/pppoe-relay
}
