# net/rp-pppoe/files/pppoe-server.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=50

DEFAULT=/etc/default/pppoe-server

start() {
	[ -f $DEFAULT ] && . $DEFAULT
	service_start /usr/sbin/pppoe-server $OPTIONS
}

stop() {
	service_stop /usr/sbin/pppoe-server
}
