# net/nut/files/nut-monitor.init
#!/bin/sh /etc/rc.common

START=60
USE_PROCD=1

restart() {
	stop_service
	start_service
}

start_service() {
	mkdir -p /var/run/nut
	chown root:root /var/run/nut
	chmod 700 /var/run/nut

	upsmon -p
}

stop_service() {
	upsmon -c stop
}

reload_service() {
	upsmon -c reload
}
