# mail/msmtp-scripts/files/msmtpq-ng-mta.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2016 OpenWrt.org

START=90

boot() {
	if [ ! -d /var/spool/msmtp ]; then
		mkdir -m1777 -p /var/spool/msmtp
	else
		chmod 1777 /var/spool/msmtp
	fi

	if [ ! -d /var/lock/msmtp ]; then
		mkdir -m1777 -p /var/lock/msmtp
	else
		chmod 1777 /var/spool/msmtp
	fi
}

