# admin/fluent-bit/files/fluent-bit.init
#!/bin/sh /etc/rc.common
START=99

USE_PROCD=1

start_service() {
	procd_open_instance
	procd_set_param command fluent-bit --supervisor -c /etc/fluent-bit/fluent-bit.yaml
	procd_set_param file /etc/fluent-bit/fluent-bit.yaml
	procd_set_param respawn
	procd_close_instance
}
