# routing/oonf-dlep-radio/files/dlep_radio.init
#!/bin/sh /etc/rc.common

START=82
DAEMON='dlep_radio'

. /lib/functions/oonf_init.sh
