# mcast-tools/files/pim6sd.init
#!/bin/sh /etc/rc.common

START=50
USE_PROCD=1

start_service() {
	procd_open_instance

	procd_set_param command /usr/sbin/pim6sd
	procd_append_param command -f -c /etc/pim6sd.conf

	procd_set_param respawn
	procd_close_instance
}
