# net/libreswan/files/etc/hotplug.d/libreswan/01-user
#!/bin/sh

[ -e "/etc/ipsec.user" ] && {
	. /etc/ipsec.user
}

exit 0
