# net/geoipupdate/files/geoipupdate.init
#!/bin/sh /etc/rc.common

START=99

start() {
	[ -d /var/GeoIP ] || {
		mkdir -p /var/GeoIP
	}
}
