# files/olsrd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

DEFAULT=/etc/default/olsrd
[ -f $DEFAULT ] && . $DEFAULT

start() {
	olsrd $OPTIONS
}

stop() {
	killall olsrd
}
