# files/olsrd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

DEFAULT=/etc/default/olsrd

start() {
	[ -f $DEFAULT ] && . $DEFAULT
	olsrd $OPTIONS
}

stop() {
	killall olsrd
}
