# net/wshaper/files/wshaper.init
#!/bin/sh /etc/rc.common

START=45
start() {
	wshaper.htb
}

stop() {
	wshaper.htb stop
}
