# nodogsplash/files/etc/uci-defaults/40_nodogsplash
#!/bin/sh

uci -q batch <<-EOF
	delete firewall.nodogsplash
	set firewall.nodogsplash=include
	set firewall.nodogsplash.type=script
	set firewall.nodogsplash.path=/usr/lib/nodogsplash/restart.sh
	commit firewall
EOF
