# files/etc/hotplug.d/iface/99-batman-adv
#!/bin/sh

. /lib/batman-adv/config.sh

bat_load_module
config_load batman-adv

case "$ACTION" in
        ifup)
                [ -d /sys/class/net/$DEVICE/batman_adv/ ] && config_foreach bat_add_interface mesh "$DEVICE"
                ;;
esac
