# bird/files/bird6.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2010-2011 OpenWrt.org

START=99
STOP=10

SERVICE_DAEMONIZE=1
SERVICE_WRITE_PID=1

BIRD="/usr/sbin/bird6 -c /etc/bird6.conf"

start() {
	service_start $BIRD -d
#	( SERVICE_MATCH_NAME=1 service_start /usr/sbin/bird6loop )
}

stop() {
#	( SERVICE_MATCH_NAME=1 service_stop /usr/sbin/bird6loop )
	service_stop $BIRD
}

reload() {
	service_reload $BIRD
}
