# net/openvswitch/files/etc/init.d/openvswitch.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2013 Julius Schulz-Zander <julius@net.t-labs.tu-berlin.de>
# Copyright (C) 2014-2017 OpenWrt.org

START=15

start() {
	/usr/share/openvswitch/scripts/ovs-ctl start
}

stop() {
	/usr/share/openvswitch/scripts/ovs-ctl stop
}

restart() {
	/usr/share/openvswitch/scripts/ovs-ctl restart
}

status() {
	/usr/share/openvswitch/scripts/ovs-ctl status
}

