# multimedia/motion/files/motion.init
#!/bin/sh /etc/rc.common

START=98
USE_PROCD=1

CONF="/etc/motion.conf"
PROG="/usr/bin/motion"

start_service() {
	procd_open_instance

	procd_set_param command "$PROG"
	procd_append_param command -c "$CONF"

	procd_close_instance
}
