# net/wifischedule/net/etc/init.d/wifi_schedule
#!/bin/sh /etc/rc.common
# Startup Script for wifi_schedule

START=99

start() {
    /usr/bin/wifi_schedule.sh startup
}
