# batmand/files/etc/init.d/vis
#!/bin/sh /etc/rc.common
START=90

start () {
	interface=$(uci get vis.general.interface)
	if [ "$interface" = "" ]; then
		echo $1 Error, you must specify at least a network interface
		exit
	fi
	vis_args=$interface

	vis $vis_args >/dev/null 2>&1
}

stop () {
        killall vis
}

        
